# homeroot/.config/autostart/mailcatcher.desktop
[Desktop Entry]
Name=mailcatcher
Exec="/usr/local/bin/mailcatcher"
Type=Application
X-GNOME-Autostart-enabled=true
